# awstats/templates/default/awstats.cron.erb
#!/bin/bash
# cron job@daily for website <%= @name %>
# generated by CHEF, do NOT change
/usr/lib/cgi-bin/awstats.pl -config=<%= @name %> -update >/dev/null
